# templates/default/cron.d/zfs-auto-snapshot
#!/bin/bash

exec <%= @binary %> --quiet --syslog --label=<%= @label %> --keep=<%= @keep_count %> <%= @datasets %>
